# local-aur.timer
[Unit]
Description=Updates the local-aur registry every 10 minutes.

[Timer]
OnBootSec=0s
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target