# local-aur.service
[Unit]
Description=Update all packages
[Service]
Type=oneshot
ExecStart=/usr/bin/local-aur build -u
